# deploy/supervisor/ozon-sheets-watch.conf
[program:ozon-sheets-watch]
command=/root/ozon-ai-agent/.venv/bin/python -m ozon_agent.cli sheets watch --interval 30
directory=/root/ozon-ai-agent
user=root
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
redirect_stderr=true
stdout_logfile=/root/ozon-ai-agent/logs/ozon-sheets-watch.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=
    PYTHONUNBUFFERED="1",
    SHEETS_DATA_SOURCE="files",
    GOOGLE_SERVICE_ACCOUNT_JSON="/root/ozon-ai-agent/secrets/google-service-account.json",
    GOOGLE_SHEETS_SPREADSHEET_ID="1MoakuEmSMkEEKf1TtoNkEF6wVyx_NBgheb3v-M4LqeI",
    SHEETS_SYNC_DELAY_SECONDS="30",
    SHEETS_RETRY_ATTEMPTS="3",
    SHEETS_RETRY_BACKOFF_SECONDS="30"
