# deploy/supervisor/ozon-telegram-bot.conf
[program:ozon-telegram-bot]
command=/root/ozon-ai-agent/.venv/bin/python -m ozon_agent.cli telegram run
directory=/root/ozon-ai-agent
user=root
autostart=false
autorestart=true
startsecs=10
stopwaitsecs=30
redirect_stderr=true
stdout_logfile=/root/.pm2/logs/ozon-telegram-bot.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=
    TELEGRAM_BOT_TOKEN="%(ENV_TELEGRAM_BOT_TOKEN)s",
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    GOOGLE_SERVICE_ACCOUNT_JSON="/root/ozon-ai-agent/secrets/google-service-account.json"
